# systemd/arch-optimize.timer
[Unit]
Description=Arch Optimize Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1month

[Install]
WantedBy=multi-user.target
