# systemd/check-failed.timer
[Unit]
Description=Check Failed Units Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=multi-user.target
