# systemd/check-failed.timer
[Unit]
Description=Check Failed Units Timer. Daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
