# systemd/check-failed.timer
[Unit]
Description=Check failed timer

[Timer]
OnStartupSec=7m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
