# systemd/check-failed.timer
[Unit]
Description=Check failed services

[Timer]
OnStartupSec=7m

[Install]
WantedBy=timers.target
