# systemd/daily-timer.target
[Unit]
Description=Basic Timer Target
StopWhenUnneeded=yes
