# systemd/empty-trash@.service
[Unit]
Description=Empty Trash Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/empty-trash %i

[Install]
WantedBy=daily-timer.target
