# systemd/empty-folders.timer
[Unit]
Description=Empty Given Folders Timer. Daily

[Timer]
OnCalendar=*-*-* 14:10:00
Persistent=true

[Install]
WantedBy=timers.target
