# systemd/installed-packages-list.timer
[Unit]
Description=List installed packages Timer. Daily

[Timer]
OnStartupSec=10m
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
